# build_files/usr/lib/systemd/system/set-cpu-min-speed.timer
[Unit]
Description=Set min CPU speed. TODO: figure out how to avoid this

[Timer]
Persistent=yes
OnCalendar=hourly

[Install]
WantedBy=timers.target
